# systemd/quant-dashboard.service
[Unit]
Description=Quant V3.1 Blazor Server Dashboard
After=quant-engine.service
Wants=quant-engine.service

[Service]
Type=simple
User=quant
WorkingDirectory=/home/quant/quant-v31/dashboard/QuantDashboard
Environment="DOTNET_ROOT=/usr/lib/dotnet"
Environment="ASPNETCORE_URLS=http://0.0.0.0:5000"
Environment="DOTNET_ENVIRONMENT=Production"
ExecStart=/usr/bin/dotnet run --no-build
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
