# systemd/quant-dashboard.service
[Unit]
Description=Quant V3.1 Blazor Server Dashboard
# NOTE: V3.1 레거시 엔진(quant-engine.service) 의존성 제거 (2026-06-25)
# 부팅 시 이 Wants=가 disabled 상태의 V3.1을 끌어올려 8000/50051 부활시키던 원인.
# 대시보드(V4)는 localhost:8001 엔진과 통신하므로 V3.1 의존성 불필요.

[Service]
Type=simple
User=quant
WorkingDirectory=/home/quant/quant-v31/dashboard/QuantDashboard
Environment="DOTNET_ROOT=/usr/lib/dotnet"
Environment="ASPNETCORE_URLS=http://0.0.0.0:5000"
Environment="DOTNET_ENVIRONMENT=Production"
ExecStart=/usr/bin/dotnet run --no-build
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
